# docker/ap/supervisor.d/20-php-fpm.conf
[program:phpfpm]
command=/usr/sbin/php-fpm7.4 --nodaemonize
autostart=true
autorestart=true
