# docker/db/supervisor.d/10-mariadb.conf
[program:mariadb]
command=/mariadb_start.sh
autostart=true
autorestart=true
