# docker/db/supervisor.d/20-phpmyadmin.conf
[program:phpmyadmin]
command=/phpmyadmin_start.sh
autostart=true
autorestart=true
